# supervisor/fix-ownership.conf
[program:fix-ownership]
command=/usr/local/bin/fix-ownership.sh
user=root
autostart=true
autorestart=false
startsecs=0
startretries=1
exitcodes=0
priority=1
stdout_logfile=/var/log/fix-ownership.log
stderr_logfile=/var/log/fix-ownership.err
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
