# supervisor/opencode-web.conf
[program:opencode-web]
command=/bin/bash -c "cd /home/cloud/proyectos && exec /usr/local/bin/opencode web --port 4096 --hostname 0.0.0.0"
user=cloud
directory=/home/cloud/proyectos
environment=HOME="/home/cloud",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",BROWSER="none",OPENCODE_SERVER_PASSWORD="%(ENV_OPENCODE_SERVER_PASSWORD)s",OPENCODE_CONFIG="/home/cloud/.config/opencode/opencode.json",OPENCODE_API_KEY="%(ENV_OPENCODE_API_KEY)s"
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/opencode-web.log
stderr_logfile=/var/log/opencode-web.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:cloudflared]
command=/usr/bin/cloudflared tunnel --config /home/cloud/.cloudflared/config.yml run --token %(ENV_CLOUDFLARE_TUNNEL_TOKEN)s --protocol http2
user=cloud
directory=/home/cloud
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/cloudflared.log
stderr_logfile=/var/log/cloudflared.err
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[program:healthcheck]
command=/bin/bash -c "while true; do curl -sf http://localhost:4096 > /dev/null && echo '[$(date)] OK' || echo '[$(date)] FAIL'; sleep 30; done"
autostart=true
autorestart=true
startsecs=0
stdout_logfile=/var/log/healthcheck.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/var/log/healthcheck.err
stderr_logfile_maxbytes=1MB
